# infra/backups/systemd/topopass-postgres-backup.timer
[Unit]
Description=Run TopoPass Postgres logical backup daily

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
Unit=topopass-postgres-backup.service

[Install]
WantedBy=timers.target
